# librecoach/rootfs/etc/s6-overlay/s6-rc.d/settings-server/run
#!/usr/bin/with-contenv bashio
exec python3 /opt/settings_server/server.py
