# librecoach/rootfs/etc/s6-overlay/s6-rc.d/vehicle-bridge/finish
#!/usr/bin/with-contenv bashio
# Publish offline status as a safety net if Python shutdown didn't complete
MQTT_USER=$(bashio::config 'mqtt_user')
MQTT_PASS=$(bashio::config 'mqtt_pass')

mosquitto_pub -h "core-mosquitto" -p 1883 \
  -u "$MQTT_USER" -P "$MQTT_PASS" \
  -t "librecoach/bridge/status" -m "offline" -r -q 1 2>/dev/null || true

mosquitto_pub -h "core-mosquitto" -p 1883 \
  -u "$MQTT_USER" -P "$MQTT_PASS" \
  -t "can/status" -m "offline" -r -q 1 2>/dev/null || true

exit 0
