# librecoach/rootfs/etc/s6-overlay/s6-rc.d/vehicle-bridge/run
#!/command/with-contenv bashio
exec python3 /opt/vehicle_bridge/main.py
